# archive/utils/supercomputer-cli/example_images/glxgears/supervisord.conf
[supervisord]
nodaemon=true
user=dev
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:tigervnc]
command=/usr/bin/Xtigervnc :1 -geometry 1280x1024 -depth 24 -rfbport 5901 -SecurityTypes None -AlwaysShared -AcceptKeyEvents -AcceptPointerEvents -AcceptCutText -SendCutText -localhost yes
user=dev
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/tigervnc.log
stderr_logfile=/var/log/supervisor/tigervnc.log

[program:openbox]
command=/usr/bin/openbox-session
user=dev
autostart=true
autorestart=true
environment=DISPLAY=:1
stdout_logfile=/var/log/supervisor/openbox.log
stderr_logfile=/var/log/supervisor/openbox.log

[program:glxgears]
command=bash -c "sleep 5 && glxgears"
user=dev
autostart=true
autorestart=true
environment=DISPLAY=:1
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/glxgears.log
stderr_logfile=/var/log/supervisor/glxgears.log

[program:novnc]
command=/usr/local/bin/run-novnc.sh
user=dev
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.log
startsecs=5
startretries=3
